# scripts/stockwatch.service
[Unit]
Description=StockWatch A股智能盯盘系统
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__STOCKWATCH_USER__
WorkingDirectory=__STOCKWATCH_HOME__/stockwatch
Environment="PATH=__STOCKWATCH_HOME__/stockwatch/.venv/bin"
Environment="PYTHONPATH=__STOCKWATCH_HOME__/stockwatch"
ExecStart=__STOCKWATCH_HOME__/stockwatch/.venv/bin/python __STOCKWATCH_HOME__/stockwatch/main.py daemon
Restart=on-failure
RestartSec=30
StandardOutput=append:__STOCKWATCH_HOME__/.stockwatch/logs/daemon.log
StandardError=append:__STOCKWATCH_HOME__/.stockwatch/logs/daemon.log

[Install]
WantedBy=multi-user.target
